# service/pnx-management-server.service
# Template file. Do not install this file directly.
# Run ./service/install.sh to render placeholders and install the unit.

[Unit]
Description=PnX Management Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__ENV_FILE__
Environment="NODE_ENV=production"
Environment="PORT=__PORT__"
Environment="HOSTNAME=0.0.0.0"
Environment="PATH=__NODE_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

ExecStartPre=__NPM_BIN__ run db:deploy
ExecStart=__NPM_BIN__ start

KillSignal=SIGTERM
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pnx-management-server

[Install]
WantedBy=multi-user.target
