# service/pnx-management-cron.service
# Template file. Do not install this file directly.
# Run ./service/install.sh to render placeholders and install the unit.

[Unit]
Description=PnX Management Cron Worker
After=network-online.target pnx-management-server.service
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__ENV_FILE__
Environment="NODE_ENV=production"
Environment="PATH=__NODE_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

ExecStart=__NPM_BIN__ run cron

KillSignal=SIGTERM
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pnx-management-cron

[Install]
WantedBy=multi-user.target
